# .devkit/scheduler/linux/repo-maintainer.cron
30 2 * * * C:/Users/murak/repos/devkit/.devkit/bin/repo-maintainer.sh
