# .devkit/scheduler/linux/repo-maintainer.timer
[Unit]
Description=Nightly repo maintainer timer for devkit

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
